# .github/workflows/resources/poseidon-minimal.socket
# This is a minimal socket definition to provide a systemd socket for the Poseidon e2e tests.

[Socket]
ListenStream = 7200
